# database/configuration/cron/acroread_cleaner.crontab

SHELL=/bin/bash

# m h dom mon dow command

* * * * * /bin/rm -f $HOME'/C:\nppdf32Log\debuglog.txt'

