# database/configuration/cron/bookmarks_builder.crontab

SHELL=/bin/bash

# m h dom mon dow command

# regenerate the bookmarks pages.
24 10 * * * (source $HOME/feisty_meow/scripts/core/profile.sh; bash $HOME/feisty_meow/scripts/bookmarks/create_marks.sh ) &>/tmp/zz_bookmarkingering_$USER.log 


