# database/configuration/cron/linux_config_snarf.crontab

SHELL=/bin/bash

# m h dom mon dow command

# do a backup of the config files once a week.
28 7 * * 4 (source $HOME/feisty_meow/scripts/core/launch_feisty_meow.sh; cd /z/stuffing/archives; perl $HOME/feisty_meow/scripts/archival/snarf_linux_config.pl) 2>&1 >>/tmp/zz_linux_config_snarf_$USER.log

