# database/configuration/cron/mysql_db_backup.crontab

SHELL=/bin/bash

# m h dom mon dow command

14 7 * * * bash $HOME/feisty_meow/scripts/database/backup_mysql_dbs.sh

