# database/configuration/cron/nechung_fortune.crontab

SHELL=/bin/bash

# m h dom mon dow command

* * * * * (export LIGHTWEIGHT_INIT=true; source $HOME/feisty_meow/scripts/core/launch_feisty_meow.sh; perl $HOME/feisty_meow/scripts/text/new_sig.pl) &>/tmp/zz_nechung_fortune_$USER.log

