# database/configuration/cron/swim_upstream.crontab

SHELL=/bin/bash

# this cron job updates the feisty meow repository on source forge based on
# the real one at feistymeow.org.

# m h dom mon dow command

USER=fred

# update the source code to the latest, every day at 3am or so.
14 */2 * * * (export LIGHTWEIGHT_INIT=true; source $HOME/feisty_meow/scripts/core/launch_feisty_meow.sh; bash $HOME/feisty_meow/scripts/rev_control/push_repo_upstream.sh $HOME/repos/sf-feistymeow-trunk ) &>>/tmp/zz_nechung_fortune_$USER.log


