# database/configuration/cron/uptime_report.crontab

SHELL=/bin/bash

# m h dom mon dow command

# the script writes the current time and current uptime to the log file.
0 * * * * (source $HOME/feisty_meow/scripts/core/launch_feisty_meow.sh; bash $HOME/feisty_meow/scripts/system/write_uptime_report.sh) 2>&1 >>/tmp/zz_uptime_report_$USER.log

