# database/configuration/cron/uptime_report.crontab

SHELL=/bin/bash

# m h dom mon dow command

# the script writes the current time and current uptime to the log file.
0 * * * * (source $HOME/yeti/scripts/core/profile.sh; bash $HOME/yeti/scripts/system/write_uptime_report.sh) &>/dev/null 

