# infobase/configuration/cron/backup_trac.crontab

##############

#[backup_trac]
# backs up trac every week.
28 7 * * 2 bash $FEISTY_MEOW_APEX/scripts/archival/backup_trac.sh /home/trac /z/stuffing/archives &>>/tmp/${CRONUSER}-cron-backup_trac.log

##############

