# infobase/configuration/cron/backup_trac.crontab



# backs up trac every week.
28 7 * * 2 bash $HOME/feisty_meow/scripts/archival/backup_trac.sh /home/trac




