# infobase/configuration/cron/backup_maildir.crontab

#[backup maildir]
# backs up the maildir folder for fred bi-monthly.
8 4 3,14 * * (source $HOME/feisty_meow/scripts/core/launch_feisty_meow.sh; tar -czf /z/stuffing/archives/maildir_backup_$(date_stringer).tar.gz $HOME/Maildir) &>/tmp/${CRONUSER}-cron-backup_maildir.log

