# infobase/configuration/cron/backup_photo_albums.crontab

#[backup photo albums]
# back up the pictures in coppermine gallery every month.
34 6 4 * * (source $HOME/feisty_meow/scripts/core/launch_feisty_meow.sh; tar -czf /z/stuffing/archives/photo_albums_backup_$(date_stringer).tar.gz /home/albums ; refred) &>/tmp/albums_backup_${USER}.log


