# infobase/configuration/cron/doxygen.crontab

SHELL=/bin/bash

# m h dom mon dow command

USER=fred

#[doxygen]
# update the feisty meow documentation once on each of tuesday, thursday and saturday.
20 4 * * 2,4,6 (source $HOME/feisty_meow/scripts/core/launch_feisty_meow.sh; cd $HOME/feisty_meow/documentation; make) &>>/tmp/zz_doxy_make_$USER.log


