# infobase/configuration/cron/emit_crontab.crontab

#[emit crontab]
# meta-operation to backup the crontab file using cron itself.
28 */8 * * * crontab -l >$HOME/$(hostname).$CRONUSER.crontab

