# infobase/configuration/cron/git_account_backup.crontab

SHELL=/bin/bash

# m h dom mon dow command

USER=fred

#[git account backup]
# archives the whole ~git folder with all account databases into a compressed tar file.
34 20 * * 5 tar -czf /z/stuffing/archives/git_archive_$(date +"\%Y_\%m_\%d_\%H\%M_\%S" | tr -d '/\n/').tar.gz /home/git >/home/git/archiving_log.txt 

