# infobase/configuration/cron/letsencrypt_renewal.crontab

##############

#[letsencrypt renewal]
# run the letsencrypt certificate renewal process every 4 hours.
# fixing issues in renewal by setting PATH.
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
37 */4 * * * (date ; letsencrypt renew) &>> /tmp/${CRONUSER}-cron-letsencrypt_renewal.txt

##############

