# infobase/configuration/cron/letsencrypt_renewal.crontab

#[letsencrypt renewal]
# run the letsencrypt certificate renewal process every 4 hours.
37 */4 * * * letsencrypt renew &>>/tmp/${USER}-cron-letsencrypt_renewal.log 

