# infobase/configuration/cron/mailman_attribute_fix.crontab

##############

#[mailman_attribute_fix]
# fix the mailman log files so that an obnoxious bug doesn't happen on the mailman web site.
23 5 * * * chmod -R 02775 /var/lib/mailman/logs &>/dev/null 

##############

