# infobase/configuration/cron/mailman_attribute_fix.crontab

#[mailman attribute fix]
# fix the mailman binaries.
23 5 * * * chmod -R 02775 /var/lib/mailman/logs &>/dev/null 

