# infobase/configuration/cron/opensim_restart.crontab

##############

#[opensim restart]
# ensure opensim is started if we can't detect it as running.
*/5 * * * * (source $FEISTY_MEOW_APEX/scripts/core/launch_feisty_meow.sh; bash $FEISTY_MEOW_SCRIPTS/opensim/maybe_restart_opensim.sh) &>>/tmp/${CRONUSER}-cron-opensim_restart.log 

##############

