# infobase/configuration/cron/process_check.crontab

SHELL=/bin/bash

# m h dom mon dow command

USER=fred

#[process check]
# a crontab for periodically dumping the process state.
* * * * * (export LIGHTWEIGHT_INIT=true; source $HOME/feisty_meow/scripts/core/launch_feisty_meow.sh; bash $HOME/feisty_meow/scripts/processes/show_hot_ones.sh) &>>/home/$USER/zz_hot_processes.log 


