# infobase/configuration/cron/snarf_linux_configuration.crontab

#[snarf linux configuration]
# do a backup of the linux configuration files once a week.
28 7 * * 4 (source $HOME/feisty_meow/scripts/core/launch_feisty_meow.sh; cd /z/stuffing/archives; perl $HOME/feisty_meow/scripts/archival/snarf_linux_config.pl) &>>/tmp/${CRONUSER}-cron-snarf_linux_configuration.log

