# infobase/configuration/cron/synch_and_swim_downstream.crontab

#[synch and swim downstream]
# updates local code repositories from master versions before uploading any changes to downstream
# repositories, such as sourceforge or github.  the uploading process (which we think of as sending
# those updated items swimming downstream to their dependent repositories) is done automatically in
# our makefiles.  that's why the code below relies on the top-level makefile in home.
14 3 * * * (source $FEISTY_MEOW_APEX/scripts/core/launch_feisty_meow.sh; cd $HOME ; bash $FEISTY_MEOW_SCRIPTS/rev_control/rev_update.sh ; make ) &>>/tmp/${CRONUSER}-cron-synch_and_swim_downstream.log

