# infobase/configuration/cron/tell_zooty_our_ip.crontab

#[tell zooty our ip]
# firewall host only--inform serene/zooty of our home IP address to enable back-connections.
0,14,28,42,56 * * * * (bash /home/fred/feisty_meow/scripts/security/tell_zooty_our_ip.sh serene.feistymeow.org fred fred ) &>>/tmp/zooty_tell_report.txt

