# infobase/configuration/cron/time_synch.crontab

#[time synch]

# standard to synch to local machine.

##############

# synchronize time/date with our nearest time server.
*/4 * * * * /usr/sbin/ntpdate greendragon &>>/tmp/${CRONUSER}-cron-time_synch.log

##############

# one machine on the network uses this version, to synch externally.
0 0,4,8,12,16,20 * * * /usr/bin/rdate -s time.nist.gov &>>/tmp/${CRONUSER}-cron-time_synch.log
#hmmm: is above still correct?

