# infobase/configuration/cron/update_system_time.crontab

##############

#[update_system_time]
# synchronize the time with a real ntp server.
42 0,4,8,12,16,20 * * * (source $FEISTY_MEOW_APEX/scripts/core/launch_feisty_meow.sh; bash $FEISTY_MEOW_APEX/scripts/system/update_system_time.sh ) &>/tmp/${CRONUSER}-cron-update_system_time.log

##############


