# infobase/configuration/cron/update_apt_based_linux.crontab


#[update apt based linux crontab]
# does a daily update of an apt-based linux box at pi o'clock.
14 3 * * * (source $HOME/feisty_meow/scripts/core/launch_feisty_meow.sh; bash $FEISTY_MEOW_SCRIPTS/system/naive_system_updater.sh) &>>/tmp/${USER}-cron-update_aptbased_linux.log


