# infobase/configuration/cron/web_group_perm_fix.crontab

# m h  dom mon dow   command

SHELL=/bin/bash

#[web group perm fix]
# fix the web directories for group access.
0 * * * * chmod -R g+w /var/www &>/dev/null


