# infobase/configuration/cron/whack_knotify4.crontab

SHELL=/bin/bash

# m h dom mon dow command

# whack this troublesome process that can go cpu crazy on a logout.
# so far this is the best solution to this major annoyance.
*/3 * * * * killall knotify4 &>/dev/null


